# pardus/playground/review/system/base/mlocate/files/locate.cron
#!/bin/sh
/usr/bin/updatedb > /dev/null
