# pardus/playground/ozan/any/bugzilla/files/cron.daily
#!/bin/bash
# Daily Bugzilla collectstats cron job run

cd %%bzinstallprefix%%/bugzilla
./collectstats.pl
