# pardus/tags/2009-EOL/util/admin/sysstat/files/sysstat.cron
# Pardus Crontab file for Sysstat
#
# Run system activity accounting tool every 10 minutes
*/10 * * * * root /usr/lib/sa/sa1 -d 1 1

# Generate a daily summary of process accounting at 23:53
53 23 * * * root /usr/lib/sa/sa2 -A
