# pardus/playground/review/network/analyzer/cacti/files/cacti.cron
*/5 * * * *     cacti   /usr/bin/php -q /usr/share/cacti/poller.php > /var/log/cacti/poller.log 2>&1
